# zabbix-smart-mon.cron
0 * * * * /usr/bin/zabbix-smart-mon -l /var/log/zabbix-smart-mon/zabbix-smart-mon.log
